# systemd_timers/hiawatha-monitor_fetch_webserver_logs.service
Unit]
Description=Run hiawatha-monitor_fetch_webserver_logs
RefuseManualStart=no
RefuseManualStop=yes
#OnFailure=failure-email@%i.service

[Service]
User=http
Type=oneshot
CPUShares=256
MemorySoftLimit=256
ExecStart=/srv/http/monitor/database/fetch_webserver_logs
