# systemd/hiawatha-monitor_send_reports.timer
[Unit]
Description=Run hiawatha-monitor_send_reports daily at 23:59

[Timer]
OnCalendar=daily 23:59


[Install]
WantedBy=timers.target
