# systemd/hiawatha-monitor_delete_old_logs.timer
[Unit]
Description=Run hiawatha-monitor_delete_old_logs daily

[Timer]
OnCalendar=daily


[Install]
WantedBy=timers.target
